# build/terminal/config/supervisor/supervisor-playback.conf
[program:playback]
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8,LC_LANG=en_US.UTF-8
command=/playback.py
numprocs=1