# build/terminal/config/supervisor/supervisor-gateone.conf
[program:gateone]
command=/usr/local/bin/gateone --user_dir=/opt/gateone/users --port=7443 --auth=api --dtach=false --session_timeout=1h 
