# build/ui/config/supervisor/supervisor-nginx.conf
[program:nginx]
command=/usr/sbin/nginx
username=root
autostart=true
