# build/ui/config/supervisor/supervisor-ip-check.conf
[program:ip_check]
command=/ip_check.sh
autorestart=false
startretries=0
priority=1
