# build/ui/config/supervisor/supervisor-incron.conf
[program:incron]
command=/usr/sbin/incrond -n
autorestart=false
startretries=0
priority=4